# src/supervisord-mariadb.conf
[program:mariadb]
command=/opt/start-mariadb.sh
numprocs=1
autostart=true
autorestart=true
